# files/etc/service/fcgiwrap/run
#!/bin/bash
[ -e /etc/default/docker-envs ] && source /etc/default/docker-envs
exec /usr/bin/spawn-fcgi -n -u www-data -s /var/run/fcgiwrap.sock -- /usr/sbin/fcgiwrap
